# scripts/sysinstall/99-monitor-hotplug.rules
ACTION=="change", SUBSYSTEM=="drm", ENV{HOTPLUG}=="1", RUN+="/home/alex/essentials/scripts/system/udev-monitor-hotplug/hotplug.sh"
